# images/claude-workspace/supervisord.conf
[supervisord]
nodaemon=true
user=dev
logfile=/var/log/supervisor/supervisord.log
pidfile=/tmp/supervisord.pid

[program:dockerd]
command=dockerd-rootless.sh
user=dev
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/dockerd.log
stderr_logfile=/var/log/supervisor/dockerd-err.log
environment=XDG_RUNTIME_DIR="/home/dev/.docker/run",HOME="/home/dev"
priority=10

[program:vscode]
command=code serve-web
    --host=0.0.0.0
    --port=8080
    --without-connection-token
    --accept-server-license-terms
    --server-data-dir=/home/dev/.vscode-server
    --user-data-dir=/home/dev/.vscode-server/user-data
    --telemetry-level=off
user=dev
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/vscode.log
stderr_logfile=/var/log/supervisor/vscode-err.log
environment=HOME="/home/dev"
priority=20
